# resources/amberstats.cron
#
# AmberStats
#
# In order to more effectively handle large numbers of inbound phone home ststs
# we bach process the incoming queue on a regular basis.
#
# This makes it easier to cache information like GeoIP rather than having to load
# it every API request.
#
*/15 * * * * apache php -q /usr/share/amberstats/htdocs/include/cron/queue_process.php >> /var/log/amberstats_queue_process.log

